# courservio-worker.ini
[program:courservio-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/virtual/$USER/courservio/artisan queue:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/virtual/$USER/courservio/storage/logs/queue.log
stopwaitsecs=3600
startsecs=30
